# com.user.clipboard-normalizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.clipboard-normalizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tyc/wx_clipboard/clipboard-normalizer</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/tmp/clipboard-normalizer.log</string>
</dict>
</plist>
